# service_files/amd-ras.service
[Unit]
Description=RAS startup script

[Service]
ExecStart=/usr/bin/amd-ras
SyslogIdentifier=amd-ras

[Install]
WantedBy=multi-user.target

